# support/crontab
MAILTO=dev@deltalcoud.apache.org
30 * * * * ${HOME}/bin/dev
